# vjudge.ini
[program: vjudge-app]
command = /home/vjudge/venv/bin/gunicorn -w 4 -b 127.0.0.1:3004 flask_app:app
directory = /home/vjudge
autostart = true
redirect_stderr = true
stdout_logfile = /etc/supervisor/log/vjudge-app.log