# vjudge.ini
[program: vjudge-app]
command = /home/vjudge/venv/bin/gunicorn flask_app:app -c gunicorn.conf.py
directory = /home/vjudge
autostart = true
redirect_stderr = true
stdout_logfile = /etc/supervisor/log/vjudge-app.log